# deployment/manual-setup/supervisor.conf
; =============================================================================
; Supervisor Configuration for Ulixai - Production
; =============================================================================
; Fichier: /etc/supervisor/conf.d/ulixai.conf
; Recharger: sudo supervisorctl reread && sudo supervisorctl update
; =============================================================================

; =============================================================================
; Queue Worker - Default (traitement général)
; =============================================================================
[program:ulixai-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ulixai.com/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --memory=128
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/ulixai.com/storage/logs/queue-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; =============================================================================
; Queue Worker - High Priority (paiements, notifications urgentes)
; =============================================================================
[program:ulixai-queue-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ulixai.com/artisan queue:work redis --queue=high --sleep=1 --tries=3 --max-time=3600 --memory=128
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/ulixai.com/storage/logs/queue-high.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; =============================================================================
; Queue Worker - Verification (Google Vision, documents)
; =============================================================================
[program:ulixai-queue-verification]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ulixai.com/artisan queue:work redis --queue=verification --sleep=3 --tries=2 --max-time=3600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/ulixai.com/storage/logs/queue-verification.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; =============================================================================
; Queue Worker - Emails (envoi de mails en arrière-plan)
; =============================================================================
[program:ulixai-queue-emails]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/ulixai.com/artisan queue:work redis --queue=emails --sleep=5 --tries=3 --max-time=3600 --memory=64
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/ulixai.com/storage/logs/queue-emails.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; =============================================================================
; Laravel Scheduler (cron jobs)
; =============================================================================
[program:ulixai-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while [ true ]; do php /var/www/ulixai.com/artisan schedule:run --verbose --no-interaction >> /var/www/ulixai.com/storage/logs/scheduler.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/ulixai.com/storage/logs/scheduler-supervisor.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3

; =============================================================================
; Group pour gérer tous les workers ensemble
; =============================================================================
[group:ulixai-workers]
programs=ulixai-queue-default,ulixai-queue-high,ulixai-queue-verification,ulixai-queue-emails,ulixai-scheduler
priority=999
